# config/includes.chroot_after_packages/etc/profile.d/start.sh
#!/bin/sh
# /etc/profile.d/start.sh
if [ -x /usr/local/bin/tb-menu ]; then
    exec /usr/local/bin/tb-menu
fi
